# deploy/systemd/eutherpal.service
[Unit]
Description=EutherPal Monopoly server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/EutherPal
Environment=EUTHERPAL_BIND=0.0.0.0:8791
Environment=EUTHERPAL_LLM_URL=http://127.0.0.1:11434/api/generate
Environment=EUTHERPAL_LLM_MODEL=supergemma4-26b-free:latest
ExecStart=%h/EutherPal/server/target/release/server
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
